# autostart/XX-Net.desktop
[Desktop Entry]
Name=XX-Net
Exec=/opt/XX-Net/start.sh
Type=Application
Terminal=false
X-GNOME-Autostart-enabled=true